# modules/cavalcade/templates/upstart.conf.erb
description "Cavalcade Daemon"
author "Chassis"

start on startup
stop on shutdown
respawn

chdir <%= @wproot %>
setuid www-data
script
  exec <%= @path %>/runner/bin/cavalcade
end script
